# pi-coding-agent-web-client-vps/examples/pi-agent.service
[Unit]
Description=Pi coding agent runner
After=network-online.target
Wants=network-online.target

[Service]
User=piagent
WorkingDirectory=/opt/pi-agent
EnvironmentFile=/opt/pi-agent/.env
ExecStart=/usr/local/bin/pi-agent-server --host 127.0.0.1 --port 7777
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
