# pi-web-client-vps-port-forward/examples/pi-vps-tunnel.service
[Unit]
Description=Pi to VPS reverse SSH tunnel
After=network-online.target
Wants=network-online.target

[Service]
User=pi
ExecStart=/usr/bin/autossh -M 0 -N -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes -i /home/pi/.ssh/pi_vps_tunnel -R 127.0.0.1:18080:127.0.0.1:3000 tunnel@your-vps.example
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
